# control/systemd/agent-watch.service
# systemd --user unit for the agent watcher. Runs as VIKTOR.
#
# Install (on the executor, as viktor):
#   loginctl enable-linger "$USER"        # keep it running when logged out (overnight!)
#   mkdir -p ~/.config/systemd/user
#   cp /home/agent/repos/plantegner/agent-runner/control/systemd/agent-watch.service \
#      ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now agent-watch.service
#
# Logs:  journalctl --user -u agent-watch -f
# Stop:  systemctl --user stop agent-watch     (or `touch ~/.config/agent-runner/PAUSED`)
#
# ExecStart points at the watcher in floorplanner's clone (where the runner lives
# until it's extracted); it serves every repo in the registry, not just that one.

[Unit]
Description=Agent runner watcher — drains repo boards onto auto/work
After=network-online.target

[Service]
Type=simple
ExecStart=/home/agent/repos/plantegner/agent-runner/bin/watch
Environment=WATCH_INTERVAL=120
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
